# supervisor-laravel-scheduler.conf
[program:laravel-scheduler]
process_name=%(program_name)s_%(process_num)02d
command=php /home/xbot_bot/xdashboard/artisan schedule:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=xbot_bot
numprocs=1
redirect_stderr=true
stdout_logfile=/home/xbot_bot/xdashboard/storage/logs/scheduler.log
stopwaitsecs=3600

